# com.studio.embed-pro.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.studio.embed-pro</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/adam/dev/embed-pro/.venv/bin/python</string>
        <string>/Users/adam/dev/embed-pro/embed-pro.py</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PORT</key>
        <string>8020</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/embed-pro.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/embed-pro.log</string>
    <key>WorkingDirectory</key>
    <string>/Users/adam/dev/embed-pro</string>
</dict>
</plist>
